# .config/systemd/user/agentsview-binary-restart.service
[Unit]
Description=Restart AgentsView after a binary update

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user try-restart agentsview.service
