# .config/systemd/user/agentsview-binary.path
[Unit]
Description=Restart AgentsView after a binary update

[Path]
PathChanged=%h/.local/bin/agentsview
Unit=agentsview-binary-restart.service

[Install]
WantedBy=default.target
