# .config/systemd/user/codex-remote-control.service
[Unit]
Description=Codex remote-control app server
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=simple
UMask=0077
WorkingDirectory=%h
ExecStart=%h/.config/codex/service
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s

[Install]
WantedBy=default.target
