# .config/lan-dns/dnsmasq-lan.service
[Unit]
Description=Split-horizon DNS for lan.domenza.in
After=network-online.target tailscaled.service lan-dns-firewall.service
Wants=network-online.target
Requires=lan-dns-firewall.service

[Service]
Type=simple
ExecStartPre=/usr/sbin/dnsmasq --test --conf-file=/etc/dnsmasq-lan.conf
ExecStart=/usr/sbin/dnsmasq --keep-in-foreground --conf-file=/etc/dnsmasq-lan.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3s
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
