# .config/lan-dns/lan-dns-sync.timer
[Unit]
Description=Periodically replicate split-horizon DNS records

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=5s
RandomizedDelaySec=10s
Unit=lan-dns-sync.service

[Install]
WantedBy=timers.target
