# .config/codex/systemd/codex-remote-control.service
[Unit]
Description=Codex remote-control app server
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=simple
UMask=0077
WorkingDirectory=%h
ExecStart=%h/.config/codex/service
Restart=always
RestartSec=5s
KillMode=mixed
KillSignal=SIGHUP
TimeoutStopSec=infinity

[Install]
WantedBy=default.target
